# systemd/linkedin-mcp-deadman.timer
[Unit]
Description=linkedin-mcp-pro — daily dead-man switch check
Documentation=https://github.com/horizonbymuneeb/linkedin-mcp-pro

[Timer]
# Run once a day at 09:00 UTC. Adjust to your local morning.
OnCalendar=daily
OnCalendar=*-*-* 09:00:00 UTC
# Catch up missed runs (e.g. if the host was down) on next boot.
Persistent=true
# Run shortly after boot if the timer was missed and is now in the past.
PersistentRandomizedDelaySec=0
# Don't run the check more than once per ~23h even if the calendar shifts.
AccuracySec=1min
# Drop the workload if the previous run is still going.
OnUnitActiveSec=1h
Unit=linkedin-mcp-deadman.service

[Install]
WantedBy=timers.target
